# debian/strongswan.postinst
#!/bin/sh
set -e

# Use the new strongswan.service instead of the legacy strongswan-starter
# systemctl unmask strongswan-starter.service
# systemctl enable strongswan-starter.service
systemctl unmask strongswan.service
systemctl enable strongswan.service

#DEBHELPER#

exit 0
