# debian/strongswan.prerm
#!/bin/sh
set -e

systemctl disable strongswan-starter.service
systemctl disable strongswan.service

#DEBHELPER#

exit 0
